# posts/docker/matrix-public-443-ops-tools-full-source/assets/source/openrc/matrix-bot-monitor
#!/sbin/openrc-run
name="Matrix bot monitor"
description="Check disk, Docker containers, TURN and LiveKit services"
command="/root/data/docker_data/synapse/scripts/run-matrix-bot-monitor-loop.sh"
command_background=true
pidfile="/run/matrix-bot-monitor.pid"
output_log="/var/log/matrix-bot-monitor-loop.log"
error_log="/var/log/matrix-bot-monitor-loop.log"
depend() {
  need net
  after docker
}
