# posts/docker/matrix-public-443-ops-tools-full-source/assets/source/openrc/matrix-simple-bot
#!/sbin/openrc-run
name="matrix-simple-bot"
description="Matrix simple invite/help bot"
command="/usr/bin/python3"
command_args="/root/data/docker_data/matrix-bot/matrix-simple-bot.py"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

export MATRIX_HOMESERVER="${MATRIX_HOMESERVER:-http://127.0.0.1:8008}"
export MATRIX_BOT_TOKEN_FILE="${MATRIX_BOT_TOKEN_FILE:-/root/data/docker_data/_secrets/matrix-bot/access_token}"
export MATRIX_SERVER_NAME="${MATRIX_SERVER_NAME:-jgaga.tk}"

depend() {
  need net
  after docker
}
