# posts/docker/matrix-public-443-ops-tools-full-source/assets/source/openrc/matrix-bot-sync
#!/sbin/openrc-run
name="Matrix bot sync"
description="Poll Matrix events for bot invites and weather commands"
command="/root/data/docker_data/synapse/scripts/run-matrix-bot-sync-loop.sh"
command_background=true
pidfile="/run/matrix-bot-sync.pid"
output_log="/var/log/matrix-bot-sync-loop.log"
error_log="/var/log/matrix-bot-sync-loop.log"
depend() {
  need net
  after docker
}
